# deploy/systemd/infra-bot.service
[Unit]
Description=Infra Bot Telegram polling service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/infra-bot/.venv/bin/infra-bot --config /etc/infra-bot/config.yaml run-bot
Restart=always
RestartSec=5
User=infra-bot
Group=infra-bot
WorkingDirectory=/opt/infra-bot/src
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
